# .config/systemd/user/emacs.target
[Unit]
Description=Emacs text editor daemons
Wants=emacs@terminal.service emacs@graphical.service

[Install]
WantedBy=default.target
