# .config/systemd/user/emacs.socket
[Socket]
ListenStream=%t/emacs.socket

[Install]
WantedBy=sockets.target
