# .config/systemd/user/emacs.service.d/spacemacs.conf
[Service]
Environment=SPACEMACSDIR=%E/spacemacs