# .config/systemd/user/emacs@.service.d/env.conf
[Service]
Environment=XAUTHORITY=%C/Xauthority
Environment=DOOMDIR=%E/doom
Environment=DOOMLOCALDIR=%h/.local/share/emacs
